# deploy/systemd/hestia-reminders.timer
[Unit]
Description=Hestia reminders — fire any due reminders every minute

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target
