# deploy/systemd/hestia-garden-watch.service
[Unit]
Description=Hestia — garden watch (soil moisture + forecast -> phone alert)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/hestia/brain
ExecStart=%h/.local/bin/uv run --project %h/hestia/brain python %h/hestia/brain/garden_watch.py
